# tldraw-frontend.service
[Unit]
Description=TLDraw Sandbox Frontend
After=network.target

[Service]
Type=simple
User=exedev
WorkingDirectory=/home/exedev/tldraw-sandbox
Environment=PATH=/home/exedev/.local/share/fnm/aliases/default/bin:/usr/local/bin:/usr/bin:/bin
ExecStart=/home/exedev/.local/share/fnm/aliases/default/bin/npx vite --port 8000 --host 0.0.0.0
Restart=always
RestartSec=3

[Install]
WantedBy=multi-user.target
